# roles/xwiki/templates/xwiki.service.j2
#
# Systemd unit file for XWiki
#
# This file is managed by Ansible, do not edit it manually.

[Unit]
Description=XWiki Jetty Application Server
Documentation=https://www.xwiki.org
After=network.target
RequiresMountsFor=/var/lib/xwiki/data

[Service]
WorkingDirectory=/var/lib/xwiki/

# Environment
Environment='JETTY_OPTS=-Djetty.http.host=127.0.0.1'

# Lifecycle
Type=simple
ExecStart=/bin/bash /usr/lib/xwiki-jetty/start_xwiki.sh --noninteractive
ExecStop=/bin/bash /usr/lib/xwiki-jetty/stop_xwiki.sh
Restart=on-abort

# Logging
SyslogIdentifier=xwiki

# Security
User=xwiki
Group=xwiki
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ReadOnlyPaths=/etc/xwiki
ReadWritePaths=/var/lib/xwiki/data/
ReadWritePaths=/var/lib/xwiki/tmp/
ReadWritePaths=/var/log/xwiki/

[Install]
WantedBy=multi-user.target
